# deploy/systemd/openjarvis.service
[Unit]
Description=Hope API Server
After=network.target

[Service]
Type=simple
User=hope
WorkingDirectory=/opt/hope
ExecStart=/opt/hope/.venv/bin/hope serve --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
Environment=HOME=/opt/hope

[Install]
WantedBy=multi-user.target
